# config/kub/leds.service
[Service]
Type=oneshot
ExecStart=/usr/bin/leds

[Install]
WantedBy=basic.target
