# system-config/kub.zahe.me/selfoss-update.service
[Unit]
Description=selfoss update feeds

[Service]
Type=oneshot
ExecStart=/usr/bin/php /usr/share/webapps/selfoss/update.php
User=http
