# system-config/linux-efi-update.path
[Unit]
Description=Linux image changed

[Path]
PathChanged=/boot/vmlinuz-linux

[Install]
WantedBy=basic.target
