# system-config/minidlna.service
.include /usr/lib/systemd/system/minidlna.service

[Service]
CapabilityBoundingSet=
ReadOnlyDirectories=/
ReadWriteDirectories=/var/cache/minidlna /var/run
Nice=15
IOSchedulingClass=idle
