# system-config/ntp.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/ntpd -qg

[Install]
WantedBy=multi-user.target
