# system-config/tor.service
.include /usr/lib/systemd/system/tor.service

[Service]
PrivateTmp=yes
ReadOnlyDirectories=/
ReadWriteDirectories=/var/lib/tor
