# system-config/user-session@.service
.include /usr/lib/systemd/system/user-session@.service

[Service]
TTYPath=/dev/tty1
EnvironmentFile=%h/.pam_environment
