# system-config/user-session@.service
.include /usr/lib/systemd/system/user-session@.service

[Service]
EnvironmentFile=%h/.pam_environment
