# config/ipfs-daemon.service
[Service]
ExecStart=/usr/bin/ipfs daemon

[Install]
WantedBy=basic.target
