# config/viber.service
[Unit]
Description=Viber

[Service]
ExecStart=/opt/viber/Viber

[Install]
WantedBy=default.target