# system-config/kub/selfoss-update.timer
[Timer]
OnCalendar=2014-*-* 01,06,11,16,20:30:00

[Install]
WantedBy=multi-user.target
