# system-config/kub/selfoss-update.service
[Unit]
Description=update selfoss

[Service]
Type=oneshot
ExecStart=/usr/bin/php /mnt/marbi/selfoss/update.php
User=http
