# config/ton/powercfg.service
[Unit]
Description=Tune power

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/powercfg

[Install]
WantedBy=basic.target
